# form-api/config/form-api.service
[Unit]
Description=Form Submission API for ImproveItMD
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/capitol/form-api/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="FLASK_SECRET_KEY=REPLACE_WITH_GENERATED_SECRET"
Environment="DASHBOARD_USERNAME=austin@improveitmd.com"
Environment="DASHBOARD_PASSWORD=REPLACE_WITH_YOUR_PASSWORD"

[Install]
WantedBy=multi-user.target
